# pihole-lcd-status/dist/pihole-lcd-status.service
[Unit]
After=network.target docker.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=/storage/pihole-lcd-status

[Install]
WantedBy=multi-user.target