# winterr/dist/weatherr.timer
[Unit]
Description=Get weather data periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
