# winterr/dist/weatherr.service
[Service]
Type=oneshot
EnvironmentFile=%h/winterr.env
ExecStart=%h/winterr --host ${HOST} --db ${DB} --tag ${TAG} weather --api-key $OWM_API_KEY

[Install]
WantedBy=default.target
