# winterr/dist/winterr.service
[Service]
Type=oneshot
EnvironmentFile=%h/%N.env
ExecStart=%h/%N --host ${HOST} --db ${DB} --tag ${TAG} sensor

[Install]
WantedBy=default.target